# bin/sucros.conf
start on started system-services
stop on runlevel [06]
script
    /usr/local/bin/sucros-daemon
end script
